# deploy/overlay/conf/sysstat/sysstat.cron
#collect stats every 1 minute
*/1 * * * * root /usr/lib64/sa/sa1 1 1